# mongorun.service
[Unit]
Description=Boot up mongo

[Service]
ExecStart=/bin/bash /home/SingleNodeMongoDB/start.sh

[Install]
WantedBy=multi-user.target